# catkin_ws/devel/etc/catkin/profile.d/10.roslaunch.sh
/home/husarion/catkin_ws/devel/.private/roslaunch/etc/catkin/profile.d/10.roslaunch.sh